# dnsmasq.conf
[program:dnsmasq]
command     = /usr/sbin/dnsmasq --no-daemon --log-queries --address=/%(ENV_TLD)s/%(ENV_IP)s
autostart   = true
autorestart = true
